# crontabs
* * * * * [ ~/myfortunes -nt ~/myfortunes.dat ] && strfile ~/myfortunes
* * * * * [ ~/myfortunes.music -nt ~/myfortunes.music.dat ] && strfile ~/myfortunes.music
